# archives/mcp-services/claude-mcp-updated.service
[Unit]
Description=Claude MCP (AI channel)
After=network-online.target
Wants=network-online.target

[Service]
User=melvin
Group=melvin
WorkingDirectory=/opt/claude
# v1.0.3 uses env var for TCP listener
Environment=MCP_LISTEN=tcp://100.112.235.46:6101
ExecStart=/usr/bin/claude mcp serve
Restart=on-failure

[Install]
WantedBy=multi-user.target